# root/container/config/crond/periodic/daily/disabled-job
#!/bin/bash

# This is an executed blank cron job